# crontab
20 10 * * * /home/advnum/top_stories_wsj_ft.py
